# hostfs/etc/services.d/uwsgi/finish
#!/bin/bash

source /scripts/logger.sh

log_interrupt "uwsgi has crashed..." "both"

s6-sleep 3
